# autobuild.service
[Unit]
Description=Autobuild repo service

[Service]
ExecStart=/data/repo/build_packages.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
